# gftool/config-gfarm/sysdep/default/gfsd.in
#!/bin/sh
#
# $Id$
#
## For RedHat Linux style rc.d, and start_service() in config-gfarm.sysdep
# chkconfig: 2345 95 05
#
## For RedHat Linux style rc.d
# description: Gfarm filesystem metaserver
# processname: gfmd
# config: @config_gfarm_gfarm_config@
# pidfile: @config_gfarm_run_dir@/gfsd.pid

PROG="gfsd"

CONF="@config_gfarm_gfarm_config@"
BINARY="@config_gfarm_prefix@/sbin/$PROG"
PID_FILE="@config_gfarm_run_dir@/${PROG}.pid"
OPTIONS="@config_gfarm_gfsd_option@"

if [ -n "@config_gfarm_globus_location@" ] &&
   [ -d "@config_gfarm_globus_location@" ]
then
	GLOBUS_LOCATION="@config_gfarm_globus_location@"
	export GLOBUS_LOCATION
	if [ -f "$GLOBUS_LOCATION/etc/globus-user-env.sh" ]; then
		. "$GLOBUS_LOCATION/etc/globus-user-env.sh"
	fi
fi

case $1 in 
start)
	if [ ! -x "$BINARY" ]; then
		echo -n "$BINARY is not installed"
		exit 1
	fi
	if [ ! -f "$CONF" ]; then
		echo -n "$CONF is not configured"
		exit 1
	fi
	$BINARY -P "$PID_FILE" $OPTIONS
	;;

stop)
	if [ -f "$PID_FILE" ] && kill -TERM `cat "$PID_FILE"`; then
		:
	elif [ -f "$PID_FILE" ]; then
		echo pid `cat "$PID_FILE"` is dead, but "$PID_FILE" remains
		rm -f "$PID_FILE"
		exit 1
	else
		echo "$PID_FILE does not exist"
		exit 1
	fi
	;;

restart)
	$0 stop
	$0 start
	;;

status)
	if [ -f "$PID_FILE" ] && kill -0 `cat "$PID_FILE"`; then
		:
	elif [ -f "$PID_FILE" ]; then
		echo pid `cat "$PID_FILE"` is dead, but "$PID_FILE" remains
		exit 1
	else
		echo "$PID_FILE does not exist. $PROG is not running?"
		exit 1
	fi
	;;

*)
	echo "Usage: $0 { start | stop | restart | status }"
	exit 1
	;;
esac
exit 0
