# 2.5.7/package/freebsd/pkg-plist
@comment $FreeBSD$
bin/config-gfarm
bin/config-agent
bin/config-gfsd
bin/gfarm.arch.guess
bin/gfarm_agent
bin/gfarm-pcp
bin/gfarm-prun
bin/gfarm-ptool
bin/gfcombine
bin/gfcombine_hook
bin/gfcp
bin/gfcp_hook
bin/gfdf
bin/gfdump
bin/gfexec
bin/gfexport
bin/gfgrep
bin/gfhost
bin/gfifo.sh
bin/gfimport_fixed
bin/gfimport_text
bin/gfkey
bin/gfls
bin/gfmkdir
bin/gfmpirun_p4
bin/gfps
bin/gfpwd
bin/gfq.sh
bin/gfq_commit.sh
bin/gfq_setup.sh
bin/gfrcmd
bin/gfreg
bin/gfrep
bin/gfrm
bin/gfrmdir
bin/gfront
bin/gfrsh
bin/gfrshl
bin/gfrun
bin/gfsched
bin/gfsck
bin/gfsetdir
bin/gfsplck
bin/gfssh
bin/gfsshl
bin/gfstat
bin/gfstatus
bin/gftest
bin/gfusage
bin/gfwhere
bin/gfwhoami
bin/thput-gfpio
etc/rc.d/gfmd.sh
etc/rc.d/gfsd.sh
include/gfarm/gfarm.h
include/gfarm/gfarm_config.h
include/gfarm/gfarm_error.h
include/gfarm/gfarm_metadb.h
include/gfarm/gfarm_misc.h
include/gfarm/gfarm_stringlist.h
include/gfarm/gflog.h
include/gfarm/gfs.h
include/gfarm/gfs_glob.h
include/gfarm/gfs_hook.h
lib/gfs_hook.o
lib/gfs_hook_debug.o
lib/gfs_hook_no_init.o
lib/gfs_hook_no_init_debug.o
lib/hooks_init_mpi.c
lib/libgfarm.so.0
lib/libgfarm.so
lib/libgfarm.a
lib/libgfs_hook_debug.so.0
lib/libgfs_hook_debug.so
lib/libgfs_hook_debug.a
lib/libgfs_hook_no_init_debug.so.0
lib/libgfs_hook_no_init_debug.so
lib/libgfs_hook_no_init_debug.a
lib/libgfs_hook.so.0
lib/libgfs_hook.so
lib/libgfs_hook.a
lib/libgfs_hook_no_init.so.0
lib/libgfs_hook_no_init.so
lib/libgfs_hook_no_init.a
libexec/gfrepbe_client
libexec/gfrepbe_server
man/ja/man1/gfarm_agent.1
man/ja/man1/gfcd.1
man/ja/man1/gfdf.1
man/ja/man1/gfexec.1
man/ja/man1/gfexport.1
man/ja/man1/gfgrep.1
man/ja/man1/gfhost.1
man/ja/man1/gfimport_fixed.1
man/ja/man1/gfimport_text.1
man/ja/man1/gfkey.1
man/ja/man1/gfls.1
man/ja/man1/gfmkdir.1
man/ja/man1/gfmpirun_p4.1
man/ja/man1/gfps.1
man/ja/man1/gfpwd.1
man/ja/man1/gfrcmd.1
man/ja/man1/gfreg.1
man/ja/man1/gfrep.1
man/ja/man1/gfrm.1
man/ja/man1/gfrmdir.1
man/ja/man1/gfront.1
man/ja/man1/gfrsh.1
man/ja/man1/gfrun.1
man/ja/man1/gfsched.1
man/ja/man1/gfsetdir.1
man/ja/man1/gfssh.1
man/ja/man1/gfstat.1
man/ja/man1/gfusage.1
man/ja/man1/gfwc.1
man/ja/man1/gfwhere.1
man/ja/man1/gfwhoami.1
man/ja/man3/gfarm.3
man/ja/man3/gfarm_hostlist_read.3
man/ja/man3/gfarm_import_fragment_config_read.3
man/ja/man3/gfarm_import_fragment_size_alloc.3
man/ja/man3/gfarm_initialize.3
man/ja/man3/gfarm_stringlist_add.3
man/ja/man3/gfarm_stringlist_cat.3
man/ja/man3/gfarm_stringlist_elem.3
man/ja/man3/gfarm_stringlist_free.3
man/ja/man3/gfarm_stringlist_init.3
man/ja/man3/gfarm_stringlist_length.3
man/ja/man3/gfarm_strings_free_deeply.3
man/ja/man3/gfarm_terminate.3
man/ja/man3/gfarm_url_fragments_replicate.3
man/ja/man3/gfarm_url_hosts_schedule.3
man/ja/man3/gfarm_url_program_deliver.3
man/ja/man3/gfarm_url_program_register.3
man/ja/man3/gfarm_url_section_replicate_from_to.3
man/ja/man3/gfarm_url_section_replicate_to.3
man/ja/man3/gfs_chdir.3
man/ja/man3/gfs_chmod.3
man/ja/man3/gfs_closedir.3
man/ja/man3/gfs_glob.3
man/ja/man3/gfs_glob_add.3
man/ja/man3/gfs_glob_free.3
man/ja/man3/gfs_glob_init.3
man/ja/man3/gfs_mkdir.3
man/ja/man3/gfs_opendir.3
man/ja/man3/gfs_pio_close.3
man/ja/man3/gfs_pio_create.3
man/ja/man3/gfs_pio_datasync.3
man/ja/man3/gfs_pio_eof.3
man/ja/man3/gfs_pio_error.3
man/ja/man3/gfs_pio_flush.3
man/ja/man3/gfs_pio_getc.3
man/ja/man3/gfs_pio_getline.3
man/ja/man3/gfs_pio_gets.3
man/ja/man3/gfs_pio_open.3
man/ja/man3/gfs_pio_putc.3
man/ja/man3/gfs_pio_putline.3
man/ja/man3/gfs_pio_puts.3
man/ja/man3/gfs_pio_read.3
man/ja/man3/gfs_pio_readdelim.3
man/ja/man3/gfs_pio_readline.3
man/ja/man3/gfs_pio_seek.3
man/ja/man3/gfs_pio_set_local.3
man/ja/man3/gfs_pio_set_view_index.3
man/ja/man3/gfs_pio_set_view_local.3
man/ja/man3/gfs_pio_sync.3
man/ja/man3/gfs_pio_truncate.3
man/ja/man3/gfs_pio_ungetc.3
man/ja/man3/gfs_pio_write.3
man/ja/man3/gfs_readdir.3
man/ja/man3/gfs_realpath.3
man/ja/man3/gfs_rename.3
man/ja/man3/gfs_rmdir.3
man/ja/man3/gfs_stat.3
man/ja/man3/gfs_stat_free.3
man/ja/man3/gfs_unlink.3
man/ja/man3/gfs_unlink_section.3
man/ja/man3/gfs_utimes.3
man/ja/man5/gfarm.conf.5
man/ja/man8/gfmd.8
man/ja/man8/gfsd.8
man/man1/gfarm_agent.1
man/man1/gfcd.1
man/man1/gfdf.1
man/man1/gfexec.1
man/man1/gfexport.1
man/man1/gfgrep.1
man/man1/gfhost.1
man/man1/gfimport_fixed.1
man/man1/gfimport_text.1
man/man1/gfkey.1
man/man1/gfls.1
man/man1/gfmkdir.1
man/man1/gfmpirun_p4.1
man/man1/gfps.1
man/man1/gfpwd.1
man/man1/gfrcmd.1
man/man1/gfreg.1
man/man1/gfrep.1
man/man1/gfrm.1
man/man1/gfrmdir.1
man/man1/gfront.1
man/man1/gfrsh.1
man/man1/gfrun.1
man/man1/gfsched.1
man/man1/gfsetdir.1
man/man1/gfssh.1
man/man1/gfstat.1
man/man1/gfusage.1
man/man1/gfwc.1
man/man1/gfwhere.1
man/man1/gfwhoami.1
man/man3/gfarm.3
man/man3/gfarm_initialize.3
man/man3/gfarm_strings_free_deeply.3
man/man3/gfarm_terminate.3
man/man3/gfarm_url_fragments_replicate.3
man/man3/gfarm_url_hosts_schedule.3
man/man3/gfarm_url_section_replicate_from_to.3
man/man3/gfarm_url_section_replicate_to.3
man/man3/gfs_chdir.3
man/man3/gfs_chmod.3
man/man3/gfs_closedir.3
man/man3/gfs_glob.3
man/man3/gfs_glob_add.3
man/man3/gfs_glob_free.3
man/man3/gfs_glob_init.3
man/man3/gfs_mkdir.3
man/man3/gfs_opendir.3
man/man3/gfs_pio_close.3
man/man3/gfs_pio_create.3
man/man3/gfs_pio_datasync.3
man/man3/gfs_pio_eof.3
man/man3/gfs_pio_error.3
man/man3/gfs_pio_flush.3
man/man3/gfs_pio_getc.3
man/man3/gfs_pio_getline.3
man/man3/gfs_pio_gets.3
man/man3/gfs_pio_open.3
man/man3/gfs_pio_putc.3
man/man3/gfs_pio_putline.3
man/man3/gfs_pio_puts.3
man/man3/gfs_pio_read.3
man/man3/gfs_pio_readdelim.3
man/man3/gfs_pio_readline.3
man/man3/gfs_pio_seek.3
man/man3/gfs_pio_set_local.3
man/man3/gfs_pio_set_view_index.3
man/man3/gfs_pio_set_view_local.3
man/man3/gfs_pio_sync.3
man/man3/gfs_pio_truncate.3
man/man3/gfs_pio_ungetc.3
man/man3/gfs_pio_write.3
man/man3/gfs_readdir.3
man/man3/gfs_realpath.3
man/man3/gfs_rename.3
man/man3/gfs_rmdir.3
man/man3/gfs_stat.3
man/man3/gfs_stat_free.3
man/man3/gfs_unlink.3
man/man3/gfs_unlink_section.3
man/man3/gfs_utimes.3
man/man5/gfarm.conf.5
man/man8/gfmd.8
man/man8/gfsd.8
sbin/gfmd
sbin/gfregister
sbin/gfsd
share/doc/gfarm/GUIDE.ja
share/doc/gfarm/Gfarm-FAQ.en
share/doc/gfarm/Gfarm-FAQ.ja
share/doc/gfarm/INSTALL.RPM.en
share/doc/gfarm/INSTALL.RPM.ja
share/doc/gfarm/INSTALL.en
share/doc/gfarm/INSTALL.ja
share/doc/gfarm/KNOWN_PROBLEMS.en
share/doc/gfarm/KNOWN_PROBLEMS.ja
share/doc/gfarm/LICENSE
share/doc/gfarm/OVERVIEW.en
share/doc/gfarm/OVERVIEW.ja
share/doc/gfarm/README.en
share/doc/gfarm/README.ja
share/doc/gfarm/README.hook.en
share/doc/gfarm/README.hook.ja
share/doc/gfarm/RELNOTES
share/doc/gfarm/SETUP.detail.en
share/doc/gfarm/SETUP.detail.ja
share/doc/gfarm/SETUP.en
share/doc/gfarm/SETUP.ja
share/doc/gfarm/html/en/ref/index.html
share/doc/gfarm/html/en/ref/man1/gfarm_agent.1.html
share/doc/gfarm/html/en/ref/man1/gfcd.1.html
share/doc/gfarm/html/en/ref/man1/gfdf.1.html
share/doc/gfarm/html/en/ref/man1/gfexec.1.html
share/doc/gfarm/html/en/ref/man1/gfexport.1.html
share/doc/gfarm/html/en/ref/man1/gfgrep.1.html
share/doc/gfarm/html/en/ref/man1/gfhost.1.html
share/doc/gfarm/html/en/ref/man1/gfimport_fixed.1.html
share/doc/gfarm/html/en/ref/man1/gfimport_text.1.html
share/doc/gfarm/html/en/ref/man1/gfkey.1.html
share/doc/gfarm/html/en/ref/man1/gfls.1.html
share/doc/gfarm/html/en/ref/man1/gfmkdir.1.html
share/doc/gfarm/html/en/ref/man1/gfmpirun_p4.1.html
share/doc/gfarm/html/en/ref/man1/gfps.1.html
share/doc/gfarm/html/en/ref/man1/gfpwd.1.html
share/doc/gfarm/html/en/ref/man1/gfrcmd.1.html
share/doc/gfarm/html/en/ref/man1/gfreg.1.html
share/doc/gfarm/html/en/ref/man1/gfrep.1.html
share/doc/gfarm/html/en/ref/man1/gfrm.1.html
share/doc/gfarm/html/en/ref/man1/gfrmdir.1.html
share/doc/gfarm/html/en/ref/man1/gfront.1.html
share/doc/gfarm/html/en/ref/man1/gfrun.1.html
share/doc/gfarm/html/en/ref/man1/gfsched.1.html
share/doc/gfarm/html/en/ref/man1/gfsetdir.1.html
share/doc/gfarm/html/en/ref/man1/gfstat.1.html
share/doc/gfarm/html/en/ref/man1/gfusage.1.html
share/doc/gfarm/html/en/ref/man1/gfwc.1.html
share/doc/gfarm/html/en/ref/man1/gfwhere.1.html
share/doc/gfarm/html/en/ref/man1/gfwhoami.1.html
share/doc/gfarm/html/en/ref/man3/gfarm.3.html
share/doc/gfarm/html/en/ref/man3/gfarm_initialize.3.html
share/doc/gfarm/html/en/ref/man3/gfarm_strings_free_deeply.3.html
share/doc/gfarm/html/en/ref/man3/gfarm_terminate.3.html
share/doc/gfarm/html/en/ref/man3/gfarm_url_fragments_replicate.3.html
share/doc/gfarm/html/en/ref/man3/gfarm_url_hosts_schedule.3.html
share/doc/gfarm/html/en/ref/man3/gfarm_url_section_replicate_from_to.3.html
share/doc/gfarm/html/en/ref/man3/gfarm_url_section_replicate_to.3.html
share/doc/gfarm/html/en/ref/man3/gfs_chdir.3.html
share/doc/gfarm/html/en/ref/man3/gfs_chmod.3.html
share/doc/gfarm/html/en/ref/man3/gfs_closedir.3.html
share/doc/gfarm/html/en/ref/man3/gfs_glob.3.html
share/doc/gfarm/html/en/ref/man3/gfs_glob_add.3.html
share/doc/gfarm/html/en/ref/man3/gfs_glob_free.3.html
share/doc/gfarm/html/en/ref/man3/gfs_glob_init.3.html
share/doc/gfarm/html/en/ref/man3/gfs_mkdir.3.html
share/doc/gfarm/html/en/ref/man3/gfs_opendir.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_close.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_create.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_datasync.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_eof.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_error.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_flush.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_getc.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_getline.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_gets.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_open.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_putc.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_putline.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_puts.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_read.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_readdelim.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_readline.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_seek.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_set_local.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_set_view_index.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_set_view_local.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_sync.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_truncate.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_ungetc.3.html
share/doc/gfarm/html/en/ref/man3/gfs_pio_write.3.html
share/doc/gfarm/html/en/ref/man3/gfs_readdir.3.html
share/doc/gfarm/html/en/ref/man3/gfs_realpath.3.html
share/doc/gfarm/html/en/ref/man3/gfs_rename.3.html
share/doc/gfarm/html/en/ref/man3/gfs_rmdir.3.html
share/doc/gfarm/html/en/ref/man3/gfs_stat.3.html
share/doc/gfarm/html/en/ref/man3/gfs_stat_free.3.html
share/doc/gfarm/html/en/ref/man3/gfs_unlink.3.html
share/doc/gfarm/html/en/ref/man3/gfs_unlink_section.3.html
share/doc/gfarm/html/en/ref/man3/gfs_utimes.3.html
share/doc/gfarm/html/en/ref/man5/gfarm.conf.5.html
share/doc/gfarm/html/en/ref/man8/gfmd.8.html
share/doc/gfarm/html/en/ref/man8/gfsd.8.html
share/doc/gfarm/html/en/user/index.html
share/doc/gfarm/html/en/user/nfs-gfarmfs.html
share/doc/gfarm/html/en/user/samba-gfarmfs.html
share/doc/gfarm/html/index.html
share/doc/gfarm/html/ja/ref/index.html
share/doc/gfarm/html/ja/ref/man1/gfarm_agent.1.html
share/doc/gfarm/html/ja/ref/man1/gfcd.1.html
share/doc/gfarm/html/ja/ref/man1/gfdf.1.html
share/doc/gfarm/html/ja/ref/man1/gfexec.1.html
share/doc/gfarm/html/ja/ref/man1/gfexport.1.html
share/doc/gfarm/html/ja/ref/man1/gfgrep.1.html
share/doc/gfarm/html/ja/ref/man1/gfhost.1.html
share/doc/gfarm/html/ja/ref/man1/gfimport_fixed.1.html
share/doc/gfarm/html/ja/ref/man1/gfimport_text.1.html
share/doc/gfarm/html/ja/ref/man1/gfkey.1.html
share/doc/gfarm/html/ja/ref/man1/gfls.1.html
share/doc/gfarm/html/ja/ref/man1/gfmkdir.1.html
share/doc/gfarm/html/ja/ref/man1/gfmpirun_p4.1.html
share/doc/gfarm/html/ja/ref/man1/gfps.1.html
share/doc/gfarm/html/ja/ref/man1/gfpwd.1.html
share/doc/gfarm/html/ja/ref/man1/gfrcmd.1.html
share/doc/gfarm/html/ja/ref/man1/gfreg.1.html
share/doc/gfarm/html/ja/ref/man1/gfrep.1.html
share/doc/gfarm/html/ja/ref/man1/gfrm.1.html
share/doc/gfarm/html/ja/ref/man1/gfrmdir.1.html
share/doc/gfarm/html/ja/ref/man1/gfront.1.html
share/doc/gfarm/html/ja/ref/man1/gfrun.1.html
share/doc/gfarm/html/ja/ref/man1/gfsched.1.html
share/doc/gfarm/html/ja/ref/man1/gfsetdir.1.html
share/doc/gfarm/html/ja/ref/man1/gfstat.1.html
share/doc/gfarm/html/ja/ref/man1/gfusage.1.html
share/doc/gfarm/html/ja/ref/man1/gfwc.1.html
share/doc/gfarm/html/ja/ref/man1/gfwhere.1.html
share/doc/gfarm/html/ja/ref/man1/gfwhoami.1.html
share/doc/gfarm/html/ja/ref/man3/gfarm.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_hostlist_read.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_import_fragment_config_read.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_import_fragment_size_alloc.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_initialize.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_stringlist_add.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_stringlist_cat.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_stringlist_elem.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_stringlist_free.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_stringlist_init.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_stringlist_length.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_strings_free_deeply.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_terminate.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_url_fragments_replicate.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_url_hosts_schedule.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_url_program_deliver.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_url_program_register.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_url_section_replicate_from_to.3.html
share/doc/gfarm/html/ja/ref/man3/gfarm_url_section_replicate_to.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_chdir.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_chmod.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_closedir.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_glob.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_glob_add.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_glob_free.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_glob_init.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_mkdir.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_opendir.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_close.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_create.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_datasync.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_eof.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_error.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_flush.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_getc.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_getline.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_gets.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_open.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_putc.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_putline.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_puts.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_read.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_readdelim.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_readline.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_seek.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_set_local.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_set_view_index.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_set_view_local.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_sync.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_truncate.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_ungetc.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_pio_write.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_readdir.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_realpath.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_rename.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_rmdir.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_stat.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_stat_free.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_unlink.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_unlink_section.3.html
share/doc/gfarm/html/ja/ref/man3/gfs_utimes.3.html
share/doc/gfarm/html/ja/ref/man5/gfarm.conf.5.html
share/doc/gfarm/html/ja/ref/man8/gfmd.8.html
share/doc/gfarm/html/ja/ref/man8/gfsd.8.html
share/doc/gfarm/html/ja/user/export-gfarm.html
share/doc/gfarm/html/ja/user/index.html
share/doc/gfarm/html/ja/user/nfs-gfarmfs.html
share/doc/gfarm/html/ja/user/samba-gfarmfs.html
share/doc/gfarm/html/ja/user/samba-hook.html
share/doc/gfarm/html/ja/user/smboverssh.html
share/doc/gfarm/html/pic/gfarm-logo.gif
share/gfarm/config/bdb.DB_CONFIG.in
share/gfarm/config/config-gfarm.ldap
share/gfarm/config/config-gfarm.postgresql
share/gfarm/config/config-gfarm.sysdep
share/gfarm/config/gfarm.conf-ldap.in
share/gfarm/config/gfarm.conf-postgresql.in
share/gfarm/config/gfarm.schema
share/gfarm/config/gfarm.schema-1.sql
share/gfarm/config/gfarm.schema-2.sql
share/gfarm/config/initial.ldif.in
share/gfarm/config/gfarm-pgsql.in
share/gfarm/config/gfarm-slapd.in
share/gfarm/config/gfmd.in
share/gfarm/config/gfarm_agent.in
share/gfarm/config/gfsd.in
share/gfarm/config/slapd.conf-2.0.in
share/gfarm/config/slapd.conf-2.1.in
share/java/gfront.jar
@unexec rmdir %%LOCALBASE%%/share/java 2>/dev/null || true
@dirrm share/gfarm/config
@dirrm share/gfarm
@dirrm share/doc/gfarm/html/pic
@dirrm share/doc/gfarm/html/ja/user
@dirrm share/doc/gfarm/html/ja/ref/man8
@dirrm share/doc/gfarm/html/ja/ref/man5
@dirrm share/doc/gfarm/html/ja/ref/man3
@dirrm share/doc/gfarm/html/ja/ref/man1
@dirrm share/doc/gfarm/html/ja/ref
@dirrm share/doc/gfarm/html/ja
@dirrm share/doc/gfarm/html/en/user
@dirrm share/doc/gfarm/html/en/ref/man8
@dirrm share/doc/gfarm/html/en/ref/man5
@dirrm share/doc/gfarm/html/en/ref/man3
@dirrm share/doc/gfarm/html/en/ref/man1
@dirrm share/doc/gfarm/html/en/ref
@dirrm share/doc/gfarm/html/en
@dirrm share/doc/gfarm/html
@dirrm share/doc/gfarm
@dirrm include/gfarm
