# debian/gfmd.conffiles
/etc/init.d/gfmd
/etc/gfarm-ldap/gfarm.schema
/etc/default/gfarm-slapd
