# debian/gfmd.conffiles
/etc/init.d/gfmd
/etc/default/gfarm-slapd
